# curator4-565072593.conf
[program:curator4-565072593]
command=python3 /home/curator4/curator4-565072593/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-565072593/test.err.log
stdout_logfile=/home/curator4/curator4-565072593/test.out.log